# deploy/systemd/zeta-sync.service.template
[Unit]
Description=ZETA-SYNC (user service)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
Environment=APP_ENV=prod
Environment=PYTHONUNBUFFERED=1
ExecStart=__PROJECT_DIR__/.venv/bin/python __PROJECT_DIR__/main.py
Restart=always
RestartSec=2
TimeoutStopSec=20

[Install]
WantedBy=default.target

